# supervisor.d/supervisor.ini
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx

[program:php7-fpm]
command=/usr/local/sbin/php-fpm
